# frizinak/mysql/mysql.supervisor.conf
[program:mysql]
user=[username]
command=mysqld --init-file=[init-file]
redirect_stderr=true
