# frizinak/apache/apache.supervisor.conf
[program:apache]
user=www-data
command=/usr/sbin/apache2ctl -D FOREGROUND
redirect_stderr=true
