# frizinak/apache/apache.supervisor.conf
[program:apache]
user=[username]
command=/usr/sbin/apache2ctl -D FOREGROUND
redirect_stderr=true
