# frizinak/ubuntu/sshd.supervisor.conf
[program:sshd]
command=/usr/sbin/sshd -De
redirect_stderr=true
